# .local/share/lxc/c1/rootfs/etc/init.d/startServer
#!/sbin/runscript
# $Header: $

depend() {
	need net
}

start() {
	ebegin "Starting Server."
	start-stop-daemon --background --start --exec /home/pi/startServer.sh
	eend $?
}

stop() {
	ebegin "Stopping Server."

	eend $?
}
